# infrastructure/systemd/timelapse-migrate.service
[Unit]
Description=Time-lapse Camera Database Migration
Documentation=man:systemd.service(5)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot

User=ubuntu
Group=ubuntu

WorkingDirectory=/opt/android-remote/current/server
EnvironmentFile=/etc/android-remote/server.env

Environment=PYTHONDONTWRITEBYTECODE=1
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/android-remote/current/server/src
Environment=TZ=UTC

ExecStartPre=/usr/bin/test \
    -d /opt/android-remote/current/server

ExecStartPre=/usr/bin/test \
    -r /opt/android-remote/current/server/alembic.ini

ExecStartPre=/usr/bin/test \
    -r /opt/android-remote/current/server/migrations/env.py

ExecStartPre=/usr/bin/test \
    -x /opt/android-remote/.venv/bin/alembic

ExecStart=/opt/android-remote/.venv/bin/alembic \
    -c /opt/android-remote/current/server/alembic.ini \
    upgrade head

RemainAfterExit=yes
TimeoutStartSec=300s
UMask=0027

StandardOutput=journal
StandardError=journal
SyslogIdentifier=timelapse-migrate

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=read-only
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictRealtime=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native

[Install]
WantedBy=timelapse-camera.target
